# pi-disable-hdmi.service
[Unit]
Description=Disable Raspberry Pi HDMI output

[Service]
Type=oneshot
ExecStart=/opt/vc/bin/tvservice -o

[Install]
WantedBy=multi-user.target
